# docker/worker/supervisord.prod.conf
; Production supervisord for the Render Background Worker. Same two processes as
; the dev config, but every process runs as www-data (the prod image's user) —
; the dev config's "app" user doesn't exist here.
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; Process the Redis queue. --max-time recycles the worker hourly to avoid memory creep.
[program:queue]
command=php artisan queue:work --tries=3 --timeout=90 --sleep=1 --max-time=3600
directory=/app
autostart=true
autorestart=true
stopwaitsecs=95
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

; Laravel scheduler (replaces system cron): dispatches due tasks every minute.
[program:scheduler]
command=php artisan schedule:work
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
